# private/installfiles/etc/systemd/system/mongod.service.d/override.conf
[Service]
Type=forking
PIDFile=/var/lib/mongodb/mongod.pid
