# private/installfiles/etc/systemd/system/codex.target
[Unit]
Description=Codex Blackboard

[Install]
WantedBy=multi-user.target
